# apps/core/dinomi-firstboot/debian/dinomi-firstboot.postinst
#!/bin/bash

PKGACTION=$1

if [ "$PKGACTION" == "configure" ] ; then
    systemctl enable dinomi-firstboot.service
fi
